# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx

[program:php]
command = /usr/sbin/php5-fpm -c /etc/php5/fpm -F

